# worker/claude-worker.service
[Unit]
Description=Claude Code Task Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/claude-worker
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
